# bootc/root/etc/systemd/system/bluetooth-init.service
[Unit]
Description=Bluetooth controller initialization
After=bluetooth.service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/bluetooth-init.sh
RemainAfterExit=yes

[Install]
RequiredBy=bluetooth.service
RemainAfterExit=yes
