# bootc/root/etc/systemd/system/microshift.service.d/train-manifests.conf
[Unit]
Wants=microshift-train-manifests.service
After=microshift-train-manifests.service
