# bootc/root/etc/systemd/system/microshift.service.d/train.conf
[Unit]
Wants=microshift-train-manifests.service bluetooth-init.service
After=microshift-train-manifests.service bluetooth-init.service
